# etc/NetworkManager/dispatcher.d/90-split-direct
#!/usr/bin/env bash
set -euo pipefail

ACTION="${2:-}"

case "$ACTION" in
  vpn-up)
    /usr/local/bin/split-direct-update
    ;;
esac