# etc/NetworkManager/dispatcher.d/90-split-direct
#!/usr/bin/env bash
set -euo pipefail

ACTION="${2:-}"
UUID="${CONNECTION_UUID:-}"

[[ "$ACTION" == "up" ]] || exit 0
[[ -n "$UUID" ]] || exit 0

TYPE="$(
  /usr/bin/nmcli -g connection.type connection show "$UUID" 2>/dev/null || true
)"

[[ "$TYPE" == "wireguard" ]] || exit 0

/usr/local/bin/split-direct-update