# etc/systemd/system/split-direct-update.service
[Unit]
Description=Update direct-routing IP rules for selected domains
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/split-direct-update